# deploy/javaforking.service
[Unit]
Description=Java Forking Service
After=network.target

[Service]
Type=forking
User=appuser
ExecStart=/usr/bin/java -jar /opt/myapp/myapp.jar &
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

